# systemd/ai-brief.service
[Unit]
Description=AI Daily Brief — build & deliver the daily AI one-pager

[Service]
Type=oneshot
ExecStart=%h/Core/Workspace/ClaudeCode/ai-daily-brief/run-brief.sh
TimeoutStartSec=1800
